# features/hello-smartcard/data/etc/udev/rules.d/99-hello-smartcard.rules
SUBSYSTEM=="usb", TAG=="security-device", ACTION=="add", RUN+="/lib/uevent/filters/smartcard-connected"
SUBSYSTEM=="usb", TAG=="security-device", ACTION=="remove", RUN+="/lib/uevent/filters/smartcard-disconnected"
