# features/hello-smartcard/data/etc/udev/rules.d/99-hello-smartcard.rules
SUBSYSTEM=="usb", ACTION=="add", TAG=="security-device", RUN+="/lib/uevent/filters/hello-smartcard"
