# features/hello-usb/data/etc/udev/rules.d/99-hello-usb.rules
SUBSYSTEM=="usb", ENV{DEVTYPE}=="usb_device", ACTION=="add", RUN+="/lib/uevent/filters/hello-usb"
